# com.commodity.tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.commodity.tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/xiangshaoxiong/Reasonix/commodityTracker/main.py</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/xiangshaoxiong/Reasonix/commodityTracker</string>

    <!-- 每天早上 9:30 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/commodity_tracker.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/commodity_tracker.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
